# com.justaboringname.studio-display-autohz.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.justaboringname.studio-display-autohz</string>
	<key>ProgramArguments</key>
	<array>
		<string>__BIN__</string>
		<string>watch</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/studio-display-autohz.out</string>
	<key>StandardErrorPath</key>
	<string>/tmp/studio-display-autohz.err</string>
</dict>
</plist>
